# packages/grappetech-neocli/grappetech-neocli-5.0.0.dev0.tar.gz/grappetech-neocli-5.0.0.dev0/neocli/config/templates/systemd/frappe-bench-frappe-schedule.service
[Unit]
Description="{{ neocli_name }}-neo-schedule"
PartOf={{ neocli_name }}-workers.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ neocli_cmd }} schedule
StandardOutput=file:{{ neocli_dir }}/logs/schedule.log
StandardError=file:{{ neocli_dir }}/logs/schedule.error.log
WorkingDirectory={{ neocli_dir }}
