# packages/grappetech-neocli/grappetech-neocli-5.0.0.dev0.tar.gz/grappetech-neocli-5.0.0.dev0/neocli/config/templates/systemd/frappe-bench.target
[Unit]
After=network.target
Requires={{ neocli_name }}-web.target {{ neocli_name }}-workers.target {{ neocli_name }}-redis.target

[Install]
WantedBy=multi-user.target
