# packages/grappetech-neocli/grappetech-neocli-5.0.0.dev0.tar.gz/grappetech-neocli-5.0.0.dev0/neocli/config/templates/systemd/frappe-bench-redis.target
[Unit]
After=network.target
Wants={{ neocli_name }}-redis-cache.service {{ neocli_name }}-redis-queue.service {{ neocli_name }}-redis-socketio.service

[Install]
WantedBy=multi-user.target
