# packages/grappetech-neocli/grappetech-neocli-5.0.0.dev0.tar.gz/grappetech-neocli-5.0.0.dev0/neocli/config/templates/systemd/frappe-bench-web.target
[Unit]
After=network.target
Wants={{ neocli_name }}-neo-web.service {{ neocli_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
