# packages/grappetech-neocli/grappetech-neocli-5.0.0.dev0.tar.gz/grappetech-neocli-5.0.0.dev0/neocli/config/templates/systemd/frappe-bench-redis-socketio.service
[Unit]
Description="{{ neocli_name }}-redis-socketio"
PartOf={{ neocli_name }}-redis.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ redis_server }} {{ redis_socketio_config }}
StandardOutput=file:{{ neocli_dir }}/logs/redis-socketio.log
StandardError=file:{{ neocli_dir }}/logs/redis-socketio.error.log
WorkingDirectory={{ sites_dir }}
