# packages/grappetech-neocli/grappetech-neocli-5.0.0.dev0.tar.gz/grappetech-neocli-5.0.0.dev0/neocli/config/templates/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:{{ neocli_name }}-neo-web]
command={{ neocli_dir }}/env/bin/gunicorn -b 127.0.0.1:{{ webserver_port }} -w {{ gunicorn_workers }} -t {{ http_timeout }} neo.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/web.log
stderr_logfile={{ neocli_dir }}/logs/web.error.log
user={{ user }}
directory={{ sites_dir }}

{% if use_rq %}
[program:{{ neocli_name }}-neo-schedule]
command={{ neocli_cmd }} schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/schedule.log
stderr_logfile={{ neocli_dir }}/logs/schedule.error.log
user={{ user }}
directory={{ neocli_dir }}

[program:{{ neocli_name }}-neo-default-worker]
command={{ neocli_cmd }} worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/worker.log
stderr_logfile={{ neocli_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ neocli_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ neocli_name }}-neo-short-worker]
command={{ neocli_cmd }} worker --queue short
priority=4
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/worker.log
stderr_logfile={{ neocli_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=360
directory={{ neocli_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ neocli_name }}-neo-long-worker]
command={{ neocli_cmd }} worker --queue long
priority=4
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/worker.log
stderr_logfile={{ neocli_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ neocli_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

{% for worker_name, worker_details in workers.items() %}
[program:{{ neocli_name }}-neo-{{ worker_name }}-worker]
command={{ neocli_cmd }} worker --queue {{ worker_name }}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/worker.log
stderr_logfile={{ neocli_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs={{ worker_details["timeout"] }}
directory={{ neocli_dir }}
killasgroup=true
numprocs={{ worker_details["background_workers"] or background_workers }}
process_name=%(program_name)s-%(process_num)d
{% endfor %}

{% else %}
[program:{{ neocli_name }}-neo-workerbeat]
command={{ neocli_dir }}/env/bin/python -m neo.celery_app beat -s beat.schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/workerbeat.log
stderr_logfile={{ neocli_dir }}/logs/workerbeat.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ neocli_name }}-neo-worker]
command={{ neocli_dir }}/env/bin/python -m neo.celery_app worker -n jobs@%%h -Ofair --soft-time-limit 360 --time-limit 390 --loglevel INFO
priority=4
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/worker.log
stderr_logfile={{ neocli_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=400
directory={{ sites_dir }}
killasgroup=true

[program:{{ neocli_name }}-neo-longjob-worker]
command={{ neocli_dir }}/env/bin/python -m neo.celery_app worker -n longjobs@%%h -Ofair --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/worker.log
stderr_logfile={{ neocli_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1540
directory={{ sites_dir }}
killasgroup=true

[program:{{ neocli_name }}-neo-async-worker]
command={{ neocli_dir }}/env/bin/python -m neo.celery_app worker -n async@%%h -Ofair --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/worker.log
stderr_logfile={{ neocli_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1540
directory={{ sites_dir }}
killasgroup=true

{% endif %}

{% if not skip_redis %}
[program:{{ neocli_name }}-redis-cache]
command={{ redis_server }} {{ redis_cache_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/redis-cache.log
stderr_logfile={{ neocli_dir }}/logs/redis-cache.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ neocli_name }}-redis-queue]
command={{ redis_server }} {{ redis_queue_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/redis-queue.log
stderr_logfile={{ neocli_dir }}/logs/redis-queue.error.log
user={{ user }}
directory={{ sites_dir }}
{% endif %}

{% if not skip_redis %}
[program:{{ neocli_name }}-redis-socketio]
command={{ redis_server }} {{ redis_socketio_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/redis-socketio.log
stderr_logfile={{ neocli_dir }}/logs/redis-socketio.error.log
user={{ user }}
directory={{ sites_dir }}
{% endif %}

{% if node %}
[program:{{ neocli_name }}-node-socketio]
command={{ node }} {{ neocli_dir }}/apps/neo/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile={{ neocli_dir }}/logs/node-socketio.log
stderr_logfile={{ neocli_dir }}/logs/node-socketio.error.log
user={{ user }}
directory={{ neocli_dir }}
{% endif %}

[group:{{ neocli_name }}-web]
programs={{ neocli_name }}-neo-web {%- if node -%} ,{{ neocli_name }}-node-socketio {%- endif%}

{% if use_rq %}

[group:{{ neocli_name }}-workers]
programs={{ neocli_name }}-neo-schedule,{{ neocli_name }}-neo-default-worker,{{ neocli_name }}-neo-short-worker,{{ neocli_name }}-neo-long-worker

{% else %}

[group:{{ neocli_name }}-workers]
programs={{ neocli_name }}-neo-workerbeat,{{ neocli_name }}-neo-worker,{{ neocli_name }}-neo-longjob-worker,{{ neocli_name }}-neo-async-worker

{% endif %}

{% if not skip_redis %}
[group:{{ neocli_name }}-redis]
programs={{ neocli_name }}-redis-cache,{{ neocli_name }}-redis-queue,{{ neocli_name }}-redis-socketio
{% endif %}
